# deploy/good-badminton-healthcheck.service
[Unit]
Description=Good Badminton health check
Wants=network-online.target
After=network-online.target good-badminton-backend.service cloudflared.service

[Service]
Type=oneshot
User=john
Group=john
WorkingDirectory=/data/projects/Good-Badminton-1
ExecStart=/data/projects/Good-Badminton-1/deploy/check_health.sh
